# ca04/systemd/system/nginx.service.d/override.conf
[Unit]
After=zerotier-one.service sys-subsystem-net-devices-ztbtosre7x.device

[Service]
ExecStartPre=/usr/bin/bash -c 'while ! ip a show dev ztbtosre7x|grep 172.23.6; do :; done'
