# ca04/systemd/system/nginx.service.d/override.conf
[Unit]
After=zerotier-one.service sys-subsystem-net-devices-ztbtosre7x.device
