# shared/systemd/system/dn42-roa.service
[Unit]
Description=Update DN42 ROA
After=network-online.target bird.service

[Service]
Type=oneshot
ExecStart=/usr/bin/update-roa.sh
